# web_dashboard/llm-benchmark-dashboard.service
[Unit]
Description=LLM Benchmark Dashboard (FastAPI)
After=network.target

[Service]
Type=simple
User=go
Group=go
WorkingDirectory=/home/go/Workspace/llm_rs2
ExecStart=/home/go/Workspace/llm_rs2/.venv/bin/python -m web_dashboard.app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
